# sh/cron/crontab
30 5 * * 1,4 /usr/local/bin/opensim_shutdown_tmux.sh > /dev/null 2>&1
